# roles/rsnapshot/templates/systemd/system/rsnapshot.timer.j2
[Unit]
Description=rsnapshot timer level {{ level }}

[Timer]
OnCalendar={{ oncalendar }}

[Install]
WantedBy=timers.target
