# roles/rsnapshot/templates/systemd/system/rsnapshot.service.j2
[Unit]
Description=rsnapshot backup

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/bin/rsnapshot {{ level }}
